# systemd/herm-backup.timer
[Unit]
Description=herm — daily backup timer (03:00 UTC)

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
