# systemd/hermes-agent.service
[Unit]
Description=Hermes Agent (herm v0.1)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=herm
Group=herm
WorkingDirectory=/home/herm
EnvironmentFile=/home/herm/.hermes/.env
ExecStart=/home/herm/.npm-global/bin/hermes-agent serve
Restart=on-failure
RestartSec=10
WatchdogSec=1800

# Hardening:
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/herm
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @mount
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
